# src/daemon/systemd/user/pipewire.socket
[Unit]
Description=PipeWire Multimedia System Socket

[Socket]
Priority=6
ListenStream=%t/pipewire-0

[Install]
WantedBy=sockets.target
